# backups/mysite.conf
# /etc/supervisor/conf.d
[program:mysite]
command = /home/max/yenvs/dj4/bin/gunicorn mysite.wsgi -b 127.0.0.1:8000 -w 5
directory = /home/max/projects/mysite
